# fle/eval/inspect/sandbox/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:factorio]
priority=1
command=/bin/sh -c "rm -rf /opt/factorio/data/elevated-rails /opt/factorio/data/quality /opt/factorio/data/space-age && exec $FACTORIO_BIN --start-server-load-scenario $FLE_SCENARIO --port 34197 --rcon-port 27015 --rcon-password factorio --server-settings /opt/factorio/config/server-settings.json --map-gen-settings /opt/factorio/config/map-gen-settings.json --map-settings /opt/factorio/config/map-settings.json --server-adminlist /opt/factorio/config/server-adminlist.json --server-banlist /opt/factorio/config/server-banlist.json --server-whitelist /opt/factorio/config/server-whitelist.json --use-server-whitelist --mod-directory /opt/factorio/mods"
environment=FACTORIO_BIN="%(ENV_FACTORIO_BIN)s",FLE_SCENARIO="%(ENV_FLE_SCENARIO)s"
autostart=true
autorestart=false
stdout_logfile=/tmp/factorio.log
stderr_logfile=/tmp/factorio.err
startsecs=5

[program:fle_bridge]
priority=2
command=/opt/fle-venv/bin/python3 /opt/fle/bridge_service.py
environment=FLE_ENV_ID="%(ENV_FLE_ENV_ID)s",FLE_SCENARIO="%(ENV_FLE_SCENARIO)s",FLE_NUM_AGENTS="%(ENV_FLE_NUM_AGENTS)s",FLE_MODS_DIR="%(ENV_FLE_MODS_DIR)s",FLE_TOOLS_DIR="%(ENV_FLE_TOOLS_DIR)s"
autostart=true
autorestart=true
stdout_logfile=/tmp/fle_bridge.log
stderr_logfile=/tmp/fle_bridge.err
startsecs=5
startretries=10
